# deployment/configs/fleet-policies/auditd-policy.json
{
  "name": "Purple-Team-Auditd",
  "description": "Linux auditd integration for kernel-level audit logging - complements Elastic Defend",
  "namespace": "default",
  "policy_id": "purple-team-linux",
  "enabled": true,
  "inputs": [
    {
      "type": "audit/auditd",
      "enabled": true,
      "streams": [
        {
          "enabled": true,
          "data_stream": {
            "type": "logs",
            "dataset": "auditd.log"
          },
          "vars": {
            "paths": ["/var/log/audit/audit.log"],
            "tags": ["caldera", "purple-team", "auditd"],
            "preserve_original_event": true
          }
        }
      ],
      "config": {
        "audit_rules": {
          "value": [
            "# ============================================================",
            "# Purple Team Audit Rules - ATT&CK Technique Detection",
            "# ============================================================",
            "",
            "# ---------------------------------------------------------",
            "# T1003 - Credential Dumping",
            "# ---------------------------------------------------------",
            "-w /etc/shadow -p wa -k credential_access",
            "-w /etc/passwd -p wa -k credential_access",
            "-w /etc/gshadow -p wa -k credential_access",
            "-w /etc/security/opasswd -p wa -k credential_access",
            "-a always,exit -F arch=b64 -S ptrace -k credential_access",
            "",
            "# ---------------------------------------------------------",
            "# T1059 - Command and Scripting Interpreter",
            "# ---------------------------------------------------------",
            "-w /bin/bash -p x -k command_execution",
            "-w /bin/sh -p x -k command_execution",
            "-w /bin/zsh -p x -k command_execution",
            "-w /usr/bin/python -p x -k command_execution",
            "-w /usr/bin/python3 -p x -k command_execution",
            "-w /usr/bin/perl -p x -k command_execution",
            "-w /usr/bin/ruby -p x -k command_execution",
            "",
            "# ---------------------------------------------------------",
            "# T1548 - Abuse Elevation Control Mechanism",
            "# ---------------------------------------------------------",
            "-w /etc/sudoers -p wa -k privilege_escalation",
            "-w /etc/sudoers.d -p wa -k privilege_escalation",
            "-w /usr/bin/sudo -p x -k privilege_escalation",
            "-w /usr/bin/su -p x -k privilege_escalation",
            "-a always,exit -F arch=b64 -S setuid -S setgid -S setreuid -S setregid -k privilege_escalation",
            "",
            "# ---------------------------------------------------------",
            "# T1547 - Boot or Logon Autostart Execution",
            "# ---------------------------------------------------------",
            "-w /etc/rc.local -p wa -k persistence",
            "-w /etc/rc.d -p wa -k persistence",
            "-w /etc/init.d -p wa -k persistence",
            "-w /etc/cron.d -p wa -k persistence",
            "-w /etc/crontab -p wa -k persistence",
            "-w /var/spool/cron -p wa -k persistence",
            "-w /etc/systemd/system -p wa -k persistence",
            "-w /usr/lib/systemd/system -p wa -k persistence",
            "",
            "# ---------------------------------------------------------",
            "# T1078 - Valid Accounts",
            "# ---------------------------------------------------------",
            "-w /var/log/lastlog -p wa -k account_access",
            "-w /var/log/wtmp -p wa -k account_access",
            "-w /var/log/btmp -p wa -k account_access",
            "-w /var/run/utmp -p wa -k account_access",
            "-w /etc/login.defs -p wa -k account_modification",
            "-w /etc/securetty -p wa -k account_modification",
            "",
            "# ---------------------------------------------------------",
            "# T1071 - Application Layer Protocol (Network)",
            "# ---------------------------------------------------------",
            "-a always,exit -F arch=b64 -S connect -F a2=16 -k network_connection",
            "-a always,exit -F arch=b64 -S accept -F a2=16 -k network_connection",
            "-a always,exit -F arch=b64 -S socket -F a0=2 -k network_socket",
            "-a always,exit -F arch=b64 -S socket -F a0=10 -k network_socket_ipv6",
            "",
            "# ---------------------------------------------------------",
            "# T1105 - Ingress Tool Transfer",
            "# ---------------------------------------------------------",
            "-w /usr/bin/wget -p x -k file_download",
            "-w /usr/bin/curl -p x -k file_download",
            "-w /usr/bin/scp -p x -k file_transfer",
            "-w /usr/bin/sftp -p x -k file_transfer",
            "-w /usr/bin/rsync -p x -k file_transfer",
            "",
            "# ---------------------------------------------------------",
            "# T1055 - Process Injection",
            "# ---------------------------------------------------------",
            "-a always,exit -F arch=b64 -S process_vm_writev -k process_injection",
            "-a always,exit -F arch=b64 -S process_vm_readv -k process_injection",
            "",
            "# ---------------------------------------------------------",
            "# T1082 - System Information Discovery",
            "# ---------------------------------------------------------",
            "-w /etc/hostname -p r -k system_discovery",
            "-w /etc/hosts -p r -k system_discovery",
            "-w /etc/resolv.conf -p r -k system_discovery",
            "-w /proc/version -p r -k system_discovery",
            "",
            "# ---------------------------------------------------------",
            "# T1562 - Impair Defenses",
            "# ---------------------------------------------------------",
            "-w /var/log -p wa -k defense_evasion",
            "-w /etc/audit/auditd.conf -p wa -k defense_evasion",
            "-w /etc/audit/audit.rules -p wa -k defense_evasion",
            "-w /sbin/auditctl -p x -k defense_evasion",
            "-w /sbin/auditd -p wa -k defense_evasion",
            "",
            "# ---------------------------------------------------------",
            "# System calls for execve (all command execution)",
            "# ---------------------------------------------------------",
            "-a always,exit -F arch=b64 -S execve -k process_execution",
            "-a always,exit -F arch=b32 -S execve -k process_execution"
          ]
        }
      }
    }
  ],
  "package": {
    "name": "auditd",
    "title": "Auditd Logs",
    "version": "3.12.0"
  },
  "vars": {},
  "metadata": {
    "purple_team": true,
    "caldera_integration": true,
    "attck_coverage": [
      "T1003",
      "T1059",
      "T1547",
      "T1548",
      "T1078",
      "T1071",
      "T1105",
      "T1055",
      "T1082",
      "T1562"
    ],
    "notes": "Auditd provides kernel-level syscall auditing. Complements Elastic Defend for comprehensive Linux telemetry. Rules mapped to MITRE ATT&CK techniques."
  }
}
